# .config/systemd/user/bspwm-reload.service
[Unit]
Description=Reload BSPWM

[Service]
Type=oneshot
ExecStart=/bin/bash -c "bspc wm -r"
StandardOutput=journal
